# face-verification.service
[Unit]
Description=Face Verification Service
Documentation=https://github.com/your-repo/face-verification
Requires=docker.service
After=docker.service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/kodyatemirov/face-verification
ExecStart=/usr/bin/docker-compose up -d --remove-orphans
ExecStop=/usr/bin/docker-compose down
ExecReload=/usr/bin/docker-compose restart
TimeoutStartSec=300
User=kodyatemirov
Group=kodyatemirov

# Политика перезапуска
Restart=on-failure
RestartSec=10s

# Переменные окружения
Environment=COMPOSE_PROJECT_NAME=face-verification
Environment=PYTHONPATH=/app

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=face-verification

[Install]
WantedBy=multi-user.target 